# templates/default/kubernetes_services_upstart.conf.erb
# ATTENTION!!! This file was generated by Chef.
# All manual changes will be overwritten.
#
description "Kubernetes services network route"

start on  (net-device-up   IFACE=<%= node['kubernetes']['interface'] %>)
stop  on  (net-device-down IFACE=<%= node['kubernetes']['interface'] %>)

exec /sbin/ip route replace <%= node['kubernetes']['api']['service_cluster_ip_range'] %> dev <%= node['kubernetes']['interface'] %>
pre-stop exec /sbin/ip route del <%= node['kubernetes']['api']['service_cluster_ip_range'] %>
