# templates/default/upstart.conf.erb
# ATTENTION!!! This file was generated by Chef.
# All manual changes will be overwritten.
#

description "<%= @service_description %>"

start on runlevel [345]
stop on runlevel [!345]

respawn

exec <%= @cmd %>
